# supervisord/delayed_job.conf
[program:delayed_job]
user=openproject
directory=/var/www/openproject
command=/var/www/openproject/docker/scripts/env.sh /var/www/openproject/docker/scripts/repeat.py 60 bundle exec bin/delayed_job run  --exit-on-complete
autorestart=true
